# automation/ansible/roles/setup_saltstack/templates/minion-override.service.j2
[Unit]
StartLimitIntervalSec=0
StartLimitBurst=6

[Service]
Type=simple
Restart=on-failure
RestartSec=60s
